# systemd/waydroid-session.service
[Unit]
Description=Start waydroid session

[Service]
ExecStart=/usr/bin/waydroid session start

[Install]
WantedBy=default.target
